# debian/zadroteeque.prerm
#!/bin/sh

systemctl stop zadroteeque.service
systemctl disable zadroteeque.service
